# web/deploy/algvex-frontend.service
[Unit]
Description=Algvex Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=linuxuser
Group=linuxuser
WorkingDirectory=/home/linuxuser/algvex/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/linuxuser/algvex/frontend
PrivateTmp=true

[Install]
WantedBy=multi-user.target
